# pyinstaller/src/honcho_pi/package_data/systemd/honcho-api.service
[Unit]
Description=Honcho API Server
Documentation=https://docs.honcho.dev
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
WorkingDirectory=/home/<user>/.local/lib/honcho
ExecStart=/home/<user>/.local/bin/uv run --no-dev fastapi dev src/main.py --host 0.0.0.0 --port 8333
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
EnvironmentFile=/home/<user>/.env
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=default.target
